# learning-companion-cloud/deploy/mac/com.ruyi.learning-companion-cloud.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.ruyi.learning-companion-cloud</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/ruyi/Desktop/ruyi-learning-plan/learning-companion-cloud/learning-companion-cloud/scripts/run_service_mac.sh</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/ruyi/Desktop/ruyi-learning-plan/learning-companion-cloud/learning-companion-cloud</string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>ThrottleInterval</key>
  <integer>10</integer>
  <key>StandardOutPath</key>
  <string>/Users/ruyi/Library/Logs/learning-companion-cloud/service.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/ruyi/Library/Logs/learning-companion-cloud/service.err.log</string>
</dict>
</plist>
